# launchd/com.user.cleanup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.user.cleanup</string>

    <key>ProgramArguments</key>
    <array>
      <string>/voller/pfad/zu/cleanup.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>86400</integer> <!-- alle 24 Stunden -->

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/cleanup.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/cleanup.err</string>
  </dict>
</plist>