# launchd/com.user.moveimages.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.user.moveimages</string>

    <key>ProgramArguments</key>
    <array>
      <string>/voller/pfad/zu/watcher.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/moveimages.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/moveimages.err</string>
  </dict>
</plist>